# webworkspace/webworkspace-sample/step1-gunicorn.socket
[Unit]
Description=Python and Django

[Socket]
ListenStream=/run/gunicorn.sock

[Install]
WantedBy=sockets.target
